# setup/docker/debian/supervisor/apache.conf
[supervisord]
nodaemon = true

[program:apache]
command = /usr/sbin/apache2ctl -D FOREGROUND
user = root
autostart = false
stdout_events_enabled=true
stderr_events_enabled=true
