# setup/docker/debian/supervisor/blackfire.conf
[program:blackfire]
command = /usr/bin/blackfire-agent --config=/home/docker/.blackfire.ini
user = root
autostart = false
stdout_events_enabled=true
stderr_events_enabled=true
